# update_server.conf
[program:update_server]
command=/usr/bin/python -m roomonitor.update_server
user=www-data
directory=/opt/room-monitor
